# scripts/bluetooth-init.service
[Unit]
Description=Unblock and power on Bluetooth adapter
After=bluetooth.service
Requires=bluetooth.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/bluetooth-init.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
